# build/rootfs/etc/init/commander.conf
# commander - configuration management server
#
# Start the configuration management server daemon

description	"CONfiguration MANagement Daemon/servER"
author          "Akshay Moghe (amoghe)"

start on stopped preflight-2 RESULT=ok
stop  on starting rcS

# Enable respawns, then configure to respawn 10 times within 5s intvl
# respawn
# respawn limit 10 5

# Reload using SIGHUP
# reload signal HUP

exec /bin/commander \
        --db-type=sqlite3 \
        --db-dsn=/config/commander/db.sq3 \
        --log-to=syslog
