# build/rootfs/etc/init/crashcorder.conf
# crashcorder - crash detection/recorder service
#

description	"crash detection and reporting daemon"
author          "Akshay Moghe (amoghe)"

# preflight-checks-stage-1 will write our config file
start on stopped preflight-checks-stage-1 RESULT=ok
stop  on runlevel [06]

# Enable respawns, then configure to respawn 10 times within 5s intvl
respawn
respawn limit 5 5

# Reload using SIGHUP
reload  signal HUP

# Start the service
exec /bin/crashcorder \
        --conf /etc/crashcorder/crashcorder.json
       # Support this option
       # --log-file=syslog
