# build/rootfs/etc/init/networking.conf
# networking - configure virtual network devices
#
# This task causes virtual network devices that do not have an associated
# kernel object to be started on boot.

description	"configure virtual network devices"

emits static-network-up
emits net-device-up

# Originally interfaces were brought up by events emitted by upstart-udev-bridge
# but since we run pretty late in the boot order (at this point, all fs are ready
# which also means virtual-fs are ready) so we can go ahead and start interfaces.
#
# (preflight-checks-stage-1 will have written out our config file)
start on stopped preflight-checks-stage-1 RESULT=ok

task

pre-start exec mkdir -p /run/network

script
  # --no-loopback is a misnomer. It means do not treat lo "differently",
  # which vanilla Ubuntu (14.04) systems do for some reason.
  exec ifup --no-loopback -a
end script
