# build/rootfs/etc/init/node_exporter.conf
# Run node_exporter (for prometheus stats)
#

description	"system/node stats exporter for prometheus"
author          "Akshay Moghe (amoghe)"

start on started commander
stop  on runlevel [06]

# Enable respawns, then configure to respawn 10 times within 5s intvl
respawn
respawn limit 5 5

# Launch as user
setuid prometheus
setgid prometheus

script
	fifofile="/tmp/node_exporter_log_fifo"
	mkfifo $fifofile
	( logger -t node_exporter < $fifofile & )
	exec > $fifofile
	rm $fifofile

	exec /usr/local/bin/node_exporter -web.listen-address "localhost:9100"
end script
