# build/rootfs/etc/init/preflight-checks-stage-2.conf
# mountall - Mount filesystems on boot
#
# This script runs system configuration after the network is available

description	"Preflight checks II (i.e. - config stage 2, post networking)"
author          "Akshay Moghe (amoghe)"

start on stopped networking
stop  on starting rcS

task

script
   echo "[NOOP] Finished preflight-checks-stage-2" >> /var/log/preflight.log
   #ruby /opt/commander/bin/preflight-check --db=/config/db.sq3 --log-file=/var/log/commander.log
end script
