# build/rootfs/etc/init/preflight-checks-stage-failed.conf
# mountall-shell - Recovery shell for configuration failure
#
# If preflight-checks exits with an error code (indicating error)
# this script starts a recovery shell that allows users to inspect
# state of the system.

description	"Recovery shell for configuration failure"

start on stopped (preflight-checks-stage-1 EXIT_STATUS=[!0] or
                  preflight-checks-stage-2 EXIT_STATUS=[!0])
stop  on runlevel [06]

task
console owner

script
    echo "System configuration failed!"
    echo "Launching recovery shell..."

    # TODO amoghe: Program to drop to db shell
    sleep(10)
end script

post-stop script
   umount -a || :
   exec reboot -f
end script
