# build/rootfs/etc/init/preflight-1.conf
# preflight stage 1 - setup system config (stage 1)
#
# This script sets up system configuration before the network is available.

description	"Preflight stage 1 (pre networking)"
author          "Akshay Moghe (amoghe)"

# Start when all filesystems are mounted, and all the mounted-* scripts
# are done reacting to the mounted events (and preparing them)
start on filesystem and (# stopped mounted-var and # We don't recieve this
                         stopped mounted-tmp and
                         stopped mounted-dev and
                         stopped mounted-run and
 			 stopped mounted-proc)
stop on starting rcS

task

script
   /bin/preflight \
        --db-type=sqlite3 \
        --db-dsn /config/db.sq3
        # TODO: --log-file=/var/log/preflight.log
end script
