# build/rootfs/etc/init/preflight-2.conf
# preflight stage 2 - setup system config (stage 2)
#
# This script sets up system configuration after the network is available

description	"Preflight stage 2 (post networking)"
author          "Akshay Moghe (amoghe)"

# Start after the networking task has stopped (and presumably, network is up)
start on stopped networking
stop  on starting rcS

task

exec echo "[NOOP] Finished preflight-2"
