# build/rootfs/etc/init/preflight-failed-shell.conf
# preflight-shell - Recovery shell for preflight failure
#
# If preflight exits to indicate that manual recovery is required, this
# starts the necessary shell.

description	"Recovery shell for preflight failure"

start on stopped preflight-1 RESULT=fail
stop on runlevel [06]

task
console owner

script
	echo "Preflight run failed."
	echo "Details of failure:"
	cat /var/log/upstart/preflight-1.log

	echo "Filesystems mounted currently:"
	mount

	echo "A maintenance shell will now be started."
	echo "CONTROL-D will terminate this shell and reboot the system."
	/sbin/sulogin
end script

post-stop script
	umount -a || :
	exec reboot -f
end script
