# build/rootfs/etc/init/radio.conf
# radio - notifications service
#

description	"alerts and notifications daemon"
author          "Akshay Moghe (amoghe)"

start on started commander
stop  on runlevel [06]

# Enable respawns, then configure to respawn 10 times within 5s intvl
respawn limit 5 5

# Reload using SIGHUP
reload  signal HUP

# Start the service
exec /bin/radio \
        --conf /etc/radio/radio.conf
        # Support these opts
        # --log-file=syslog
