# build/rootfs/etc/update-motd.d/00-header
#!/bin/bash

name=$(hostname)

echo ""
cat << EOF
 _ __ ___   ___| | _____| |_ ___| |__ (_)_ __
| '__/ _ \ / __| |/ / _ \ __/ __| '_ \| | '_ \\
| | | (_) | (__|   <  __/ |_\__ \ | | | | |_) |
|_|  \___/ \___|_|\_\___|\__|___/_| |_|_| .__/
                                      |_|

EOF

echo "Welcome to ${name}"
echo ""
