# build/rootfs/etc/update-motd.d/10-help-text
#!/bin/bash

echo "For help/troubleshooting, contact support at support@rocketship.com"
echo ""
